# deploy/supervisor-octane.conf
; =============================================================================
; NiceShoply Octane 应用服务进程（Supervisor 配置样例，可选）
; 对应《NiceShoply改进任务清单》§4.2
; -----------------------------------------------------------------------------
; 仅在启用 Octane（FrankenPHP / Swoole / RoadRunner）时使用。
; 安装位置（示例）：/etc/supervisor/conf.d/niceshoply-octane.conf
; 生效命令：
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start niceshoply-octane
;
; 说明：
;   - 请将 /path/to 替换为真实部署路径
;   - --workers / --task-workers 根据 CPU 核数调整
;   - 代码发布后需执行：php artisan octane:reload
; =============================================================================
[program:niceshoply-octane]
process_name=%(program_name)s
command=php /path/to/niceshoply/artisan octane:start --host=127.0.0.1 --port=8000 --workers=4 --task-workers=2
directory=/path/to/niceshoply
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/path/to/niceshoply/storage/logs/octane.log
stopwaitsecs=30
