# deploy/supervisor-horizon.conf
; =============================================================================
; NiceShoply Horizon 队列守护进程（Supervisor 配置样例）
; 对应《NiceShoply改进任务清单》§4.2
; -----------------------------------------------------------------------------
; 安装位置（示例）：/etc/supervisor/conf.d/niceshoply-horizon.conf
; 生效命令：
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start niceshoply-horizon
;
; 说明：
;   - 请将 /path/to 替换为真实部署路径（artisan 所在目录）
;   - user 改为运行 PHP-FPM/Web 的系统用户（如 www-data）
;   - Horizon 自身管理 worker 进程，因此 numprocs=1
; =============================================================================
[program:niceshoply-horizon]
process_name=%(program_name)s
command=php /path/to/niceshoply/artisan horizon
directory=/path/to/niceshoply
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/path/to/niceshoply/storage/logs/horizon.log
stopwaitsecs=3600
; Horizon 平滑退出：先发 SIGTERM 等待当前任务完成
stopsignal=TERM
